# idleServer/supervisord/nginx.conf
[program:nginx]
command=/sbin/nginx -g "daemon off;"
numprocs=1
autostart=true
autorestart=false